# rootfs/etc/services.d/smartd/run
#!/usr/bin/with-contenv bash

echo "starting smartd"
/usr/sbin/smartd --debug
